# templates/resources/service_openrc.md.tmpl
---
# generated by https://github.com/hashicorp/terraform-plugin-docs
page_title: "{{.Name}} | {{.Type}} | {{.ProviderName}}"
name: "{{.Name}}"
type: "{{.Type}}"
subcategory: ""
description: |-
{{ .Description | plainmarkdown | trimspace | prefixlines "  " }}
---

# {{.Type}}: {{.Name}}

{{ .Description | trimspace }}

## Usage

### Start a service

This example ensures that the OpenRC service `nginx` is started. The service must exist.

```terraform
resource "system_service_openrc" "nginx" {
  name    = "nginx"
  status  = "started"
}
```

### Enable a service

This example ensures that the OpenRC service `nginx` is enabled. The service is enabled for the runlevel defined in the attribute `runlevel` which defaults to `default`. The service must exist.

```terraform
resource "system_service_openrc" "nginx" {
  name    = "nginx"
  enabled = true
}
```

### Create, enable, and start a service

This example creates an OpenRC service script `busybox-httpd` using the `system_file` resource and subsequently enables and starts the service using the `system_service_openrc`. The service starts a busybox httpd daemon.

```terraform
resource "system_file" "busybox_httpd" {
  path    = "/etc/init.d/busybox-httpd"
  mode    = 755
  user    = "root"
  group   = "root"
  content = <<EOT
#!/sbin/openrc-run

name=$SVCNAME
command="/bin/busybox-extras httpd"
command_args="-p 8080 -h /var/www/html"

depend() {
    need net localmount
    after firewall
}

start_pre() {
    mkdir -p /var/www/html
}
EOT
}

resource "system_service_openrc" "busybox_httpd" {
  name    = system_file.busybox_httpd.basename
  enabled = true
  status  = "started"
}
```

## Notes

This section describes general notes for using the `system_service_openrc` resource.

- The service script must exist.
- The resource does not manage, create, or delete the service script.
- The resource remembers the `enabled` state and the `status` state at the time the resource is created. This state is referred to as the *original state*.
- When the resource is deleted, the service is reverted to the original state.
- Avoid defining multiple `system_service_openrc` resources, which manage the same service in the same Terraform configuration. Instead, merge all attributes in a single `system_service_openrc` resource.

{{ if .HasExample -}}
## Example Usage

{{ printf "{{tffile %q}}" .ExampleFile }}
{{- end }}

{{ .SchemaMarkdown | trimspace }}

{{ if .HasImport -}}
## Import

Import is supported using the following syntax:

{{ printf "{{codefile \"shell\" %q}}" .ImportFile }}
{{- end }}